# assets/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php-fpm7.2 --nodaemonize
stdout_logfile = /dev/stdout
redirect_stderr = true
autorestart = true
priority = 10
